# xdg/.config/autostart/swap_caps_lock_with_escape.desktop
[Desktop Entry]
Type=Application
Name=Set caps lock to be another escape key
Exec=setxkbmap -option caps:escape
Terminal=false
